# supervisor.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
redirect_stderr=true
